# xuidb/x-ui
#!/sbin/openrc-run

name="x-ui"

command="/usr/local/x-ui/x-ui"

directory="/usr/local/${name}"

pidfile="/var/run/${name}.pid"

command_background="yes"

depend() {
	need networking 
}

