# xuidb/x-ui
#!/sbin/openrc-run

name="x-ui"
command=/usr/local/x-ui/x-ui
directory="/usr/local/${name}"

depend() {
	need networking 
}

